# templates/default/bamboo-server.conf.erb
# upstart file at /etc/init/bamboo-server.conf

description "Bamboo"

limit nofile 32768 65536

start on runlevel [2345]
stop on [!12345]

respawn
respawn limit 2 5


umask 007

kill timeout 300

exec <%= node['qubit_bamboo']['home'] -%>/bamboo -config <%= node['qubit_bamboo']['config_path']-%> -log <%= node['qubit_bamboo']['log_path'] %> 
